# felports/keyd/keyd.initd
#!/sbin/openrc-run

command=/usr/bin/keyd
